# root/etc/init/btwilink.conf
description "WiLink BT module"

start on started uim-sysfs

script
	modprobe btwilink
end script
