# root/etc/init/uim-sysfs.conf
description "Usermode Init Manager for TI Shared Transport"

start on net-device-up IFACE=wlan0

exec /usr/local/bin/uim