# root/etc/init/uim-sysfs.conf
description "Usermode Init Manager for TI Shared Transport"

start on filesystem

exec /usr/local/bin/uim